# homeassistant/addon/homesec/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/with-contenv bashio
exec su postgres -s /bin/sh -c "postgres -D /data/postgres/data -c listen_addresses='127.0.0.1'"
